# packages/ubo-app/ubo_app-0.7.19-py3-none-any.whl/ubo_app/system_services/update.service
[Install]
WantedBy=multi-user.target

[Unit]
Description=Ubo Update Service

[Service]
Type=oneshot
User=root
Environment=INSTALLATION_PATH={{INSTALLATION_PATH}}
ExecStart={{INSTALLATION_PATH}}/env/lib/python3.11/site-packages/ubo_app/system_services/update.sh
RemainAfterExit=no
